# .config/systemd/user/exdl.path
[Unit]
Description=Extract zips and rars when /home/paul/Downloads is not empty.
PartOf=nzbget.service

[Path]
#DirectoryNotEmpty=/home/paul/Downloads
#PathExistsGlob=/home/paul/Downloads/*
PathChanged=/home/paul/Downloads

[Install]
WantedBy=paths.target
